# systemd/gbrain-http-wrapper.service
[Unit]
Description=GBrain HTTP Wrapper (MCP over HTTP+Bearer)
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=simple
WorkingDirectory=/home/ec2-user/gbrain-http-wrapper
EnvironmentFile=/home/ec2-user/gbrain-http-wrapper/.env
ExecStart=/home/ec2-user/.bun/bin/bun run src/server.ts
Restart=on-failure
RestartSec=30
TimeoutStartSec=30
StandardOutput=journal
StandardError=journal
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=default.target
